# remote_files/supervisord/panopticon.ini
[program:panopticon]
command=/etc/alternatives/java -cp /etc/panopticon:/usr/local/lib/panopticon/panopticon-standalone.jar panopticon.core --config /etc/iplant/de/panopticon.properties
directory=/                    ; directory to cwd to before exec (def no cwd)
startsecs=5                    ; number of secs prog must stay running (def. 1)
startretries=3                 ; max # of serial start failures (default 3)
exitcodes=0                    ; 'expected' exit codes for process (default 0,2)
user=condor                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/panopticon/panopticon-shell.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10      ; # of stdout logfile backups (default 10)
