# roles/user-sessions/files/user-sessions.ini
[program:user-sessions]
command=/etc/alternatives/java -cp /etc/user-sessions:/usr/local/lib/user-sessions/user-sessions-standalone.jar user_sessions.core --config /etc/iplant/de/user-sessions.properties
directory=/                    ; directory to cwd to before exec (def no cwd)
startsecs=5                    ; number of secs prog must stay running (def. 1)
startretries=3                 ; max # of serial start failures (default 3)
exitcodes=0                    ; 'expected' exit codes for process (default 0,2)
user=iplant                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/user-sessions/user-sessions-shell.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10      ; # of stdout logfile backups (default 10)
