# roles/docker-cleanup/files/docker-cleanup.timer
[Unit]
Description=Docker Cleanup Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
