# roles/info-typer/files/info-typer.ini
[program:info-typer]
command=/etc/alternatives/java -cp /etc/info-typer:/usr/local/lib/info-typer/info-typer-standalone.jar info_typer.core --config /etc/iplant/de/info-typer.properties
directory=/                                        ; directory to cwd to before exec (def no cwd)
startsecs=5                                        ; number of secs prog must stay running (def. 1)
startretries=3                                     ; max # of serial start failures (default 3)
exitcodes=0                                        ; 'expected' exit codes for process (default 0,2)
user=iplant                                        ; setuid to this UNIX account to run the program
redirect_stderr=true                               ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/info-typer/info-typer-shell.log ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB                       ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10                          ; # of stdout logfile backups (default 10)
